# update-mirror.cron
* * * * * root sh /root/update-mirror.sh
